# resources/dataprev/dev/horizon.conf
[supervisord]
nodaemon=true

[program:horizon]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan horizon --queue=default
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/horizon.log

[program:horizon-siape]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan horizon --queue=siape_queue
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/horizon-siape.log